# local/voice-assistant.service
[Unit]
Description=Local Voice Assistant Service
After=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/ai-voice-assistant-rp/local
Environment=PYTHONPATH=/home/pi/ai-voice-assistant-rp
Environment=PORCUPINE_ACCESS_KEY=your_key_here
Environment=GOOGLE_API_KEY=your_key_here
Environment=TAVILY_API_KEY=your_key_here
ExecStart=/usr/bin/python3 /home/pi/ai-voice-assistant-rp/local/simple_local_assistant.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
